# s3_backup_postgres@.timer
[Unit]
Description=Back up Postgres DB %I to S3

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true
