# s3_backup_dir@.timer
[Unit]
Description=Back up directory %I to S3

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=weekly
AccuracySec=1h
Persistent=true
