# s3_backup_mysql@.service
[Unit]
Description=MySQL DB %I backup to S3
After=network.target
Requisite=mysql.service

[Service]
Type=oneshot
ExecStart=s3_backup_mysql.sh %I
