# postgres/s3_backup_postgres@wal.timer.d/calendar.conf
[Timer]
OnCalendar=daily
Persistent=false
