# pg/backup-postgres@.service.d/aws-cli.conf
[Service]
Environment=AWS_CONFIG_FILE=/etc/aws.ini
