# pg/backup-postgres@.service.d/bucket.conf
[Service]
Environment=S3_BUCKET_NAME=mybucket
