# pg/backup-postgres@.service.d/filename.conf
[Service]
Environment=S3_FILE_NAME=%%F.tgz
