# pg/backup-postgres@.service.d/waldir.conf
[Service]
Environment=WAL_DIR=/var/spool/postgresql
