# dir/dir@.service
[Unit]
Description=Directory %I backup to S3
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/s3_backup_dir.sh %I
