# pg/backup-postgres@wal.service.d/s3-prefix.conf
[Service]
Environment=S3_PREFIX=database/wal/
